# templates/etc/systemd/system/carbon-c-relay.service.erb
[Unit]
Description=Starts and stops carbon-c-relay daemon
After=network.target

[Service]
EnvironmentFile=/etc/sysconfig/<%= @service_name %>
User=<%= @user %>
Group=<%= @group %>
PIDFile=<%= @pid_dir %>/<%= @service_name %>.pid
ExecStart=/bin/carbon-c-relay -f <%= @config_file %>
Type=forking

[Install]
WantedBy=multi-user.target