# templates/etc/systemd/system/carbon-c-relay.service.erb
[Unit]
Description=carbon-c-relay daemon
After=network.target

[Service]
EnvironmentFile=/etc/sysconfig/<%= @package_name %>
User=<%= @user %>
Group=<%= @group %>
PIDFile=<%= @pid_dir %>/<%= @package_name %>.pid
ExecStart=/bin/carbon-c-relay -f <%= @config_file %> $ARGS
Type=forking

[Install]
WantedBy=multi-user.target
