# templates/etc/systemd/system/carbon-c-relay.service.erb
# <%= @init_file %>: Managed by puppet.
[Unit]
Description=carbon-c-relay daemon
After=network.target

[Service]
EnvironmentFile=/etc/sysconfig/<%= @package_name %>
User=<%= @user %>
Group=<%= @group %>
ExecStart=/bin/carbon-c-relay -f <%= @config_file %> $ARGS

[Install]
WantedBy=multi-user.target
